# nextcloud/nextcloud-cron.timer
[Unit]
Description=CronJobs for Nextcloud

[Timer]
OnCalendar=*:0/5:00
Persistent=true
Unit=nextcloud-cron.service

[Install]
WantedBy=basic.target
