# searxng/searxng.container
[Unit]
Description=Searxng Main Container
Wants=network-online.target
After=network-online.target

[Container]
ContainerName=searxng
Image=docker.io/searxng/searxng:latest
Pod=searxng.pod
AutoUpdate=registry
Volume=$PWD/config:/etc/searxng:Z
Volume=searxng_cache:/var/cache/searxng:Z
PublishPort=9080:8080/tcp
PublishPort=9080:8080/udp

[Service]
Restart=always

[Install]
WantedBy=default.target
